# data/beluga.service
[Unit]
Description=Beluga Discord bot
After=network.target beluga.socket
Requires=beluga.socket

[Service]
Type=notify
ExecStart=/usr/bin/beluga start
ExecStop=/usr/bin/beluga stop
User=beluga
Group=beluga

[Install]
WantedBy=multi-user.target
